# install/cinnamon/files/wallpaper.timer
[Unit]
Description=Execute wallpaper switcher service every hour

[Timer]
OnCalendar=hourly
Unit=wallpaper.service

[Install]
WantedBy=timers.target
